# .config/systemd/user/nextcloud.service
[Unit]
Description=Launch NextCloud Client

[Service]
Type=simple
ExecStart=/usr/bin/nextcloud --background

[Install]
WantedBy=default.target
