# .config/systemd/user/redshift-gtk.service
[Unit]
Description=Redshift GTK

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/redshift-gtk

[Install]
WantedBy=default.target
